# Godeps/_workspace/src/github.com/Shopify/sarama/vagrant/kafka.conf
start on started zookeeper-ZK_PORT
stop on stopping zookeeper-ZK_PORT

pre-start exec sleep 2
exec /opt/kafka-KAFKAID/bin/kafka-server-start.sh /opt/kafka-KAFKAID/config/server.properties
