# Godeps/_workspace/src/github.com/Shopify/sarama/vagrant/toxiproxy.conf
start on started networking
stop on shutdown

env KAFKA_INSTALL_ROOT=/opt

exec /opt/run_toxiproxy.sh
