# Godeps/_workspace/src/github.com/Shopify/sarama/vagrant/zookeeper.conf
start on started toxiproxy
stop on stopping toxiproxy

exec /opt/kafka-KAFKAID/bin/zookeeper-server-start.sh /opt/kafka-KAFKAID/config/zookeeper.properties
